# docker/s6-rc.d/caddy/run
#!/command/execlineb -P
cd /srv
caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
