# docker/s6-rc.d/server/run
#!/command/execlineb -P
with-contenv
cd /app
./server serve
