# docker/s6-rc.d/server/run
#!/command/execlineb -P
cd /app
./server serve
